# build/debian/miru.service
[Unit]
Description=Miru Config Agent
After=network.target
Requires=miru.socket

[Service]
Type=simple

# user and group
Group=miru
User=miru

# runtime
StateDirectory=miru
LogsDirectory=miru
ExecStart=/usr/sbin/miru-agent
WorkingDirectory=/var/lib/miru

# logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=miru-agent

# security
# ReadWritePaths is kept even though ProtectSystem is not enabled by
# default. It is a no-op in the default configuration, but becomes
# meaningful if a customer ships a hardening drop-in that enables
# ProtectSystem=strict — the agent's own working directories are then
# writable without the customer needing to remember to add them.
ReadWritePaths=/var/lib/miru /var/log/miru /srv/miru
# Prevents the service and its children from gaining new privileges
NoNewPrivileges=true
# Gives the service a private /tmp and /var/tmp
PrivateTmp=true
# Restricts access to /dev, only allowing minimal device nodes
PrivateDevices=true
# Restricts the service’s view of /proc to only its own processes (Linux 5.8+)
ProtectProc=invisible
# Makes kernel tunable files (e.g., /proc/sys, /sys) read-only or inaccessible
ProtectKernelTunables=true
# Prevents loading or unloading of kernel modules
ProtectKernelModules=true
# Prevents modification of cgroup settings
ProtectControlGroups=true
# Limits network protocols to UNIX sockets, IPv4, and IPv6
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# Prevents gaining real-time scheduling privileges
RestrictRealtime=true
# Prevents creating or executing set-user-ID or set-group-ID files
RestrictSUIDSGID=true


[Install]
WantedBy=multi-user.target
